# settingsd.service
[Unit]
Description=Settings Daemon

[Service]
Type=dbus
BusName=org.os.Settings
ExecStart=settingsd

[Install]
WantedBy=multi-user.target
